# systemd/wb-rosdomofon.service
[Unit]
Description=WB Rosdomofon integration service
After=network.target mosquitto.service
Wants=mosquitto.service
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=root
WorkingDirectory=/mnt/data/wb-rosdomofon
Environment=NODE_ENV=production
ExecStartPre=/usr/bin/test -f /mnt/data/wb-rosdomofon/wb-rosdomofon.cfg
ExecStart=/usr/bin/node src/index.cjs
Restart=always
RestartSec=3
TimeoutStopSec=5
KillSignal=SIGTERM
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
